# random-background.timer
[Unit]
Description=Select random background daily (timer)

[Timer]
OnCalendar=daily
Unit=%i.service

[Install]
WantedBy=default.target
